# config/com.recoblog.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.recoblog.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source /Users/evanhwang/my-github/project/recoblog/.venv/bin/activate && python3 -m scripts.publish --date today</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/evanhwang/my-github/project/recoblog/logs/publish.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/evanhwang/my-github/project/recoblog/logs/publish.error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/evanhwang/my-github/project/recoblog</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
